# docker_base/etc/services.d/nerve/run
#!/usr/bin/with-contenv sh

exec nerve -i "`hostname`" -c /opt/smartstack/nerve/config.json
