# docker_base/etc/services.d/nervesync/run
#!/usr/bin/with-contenv sh

exec nervesync -H "$(echo $HOST_IP)" -z "$(echo $ZK_HOSTS)"
